# supervisor.ini
[program:newspaper]
directory = /root/projects/newspaper
environment = PATH="/root/projects/newspaper/ENV/bin"
command = /root/projects/newspaper/ENV/bin/python AutoGreen.py --browser chrome
process_name = newspaper-%(process_num)s
autorestart = true
redirect_stderr = true
stdout_logfile = /var/log/newspaper-supervisor-stdout.log
stdout_logfile_maxbytes = 100MB
stdout_logfile_backups = 50
stdout_capture_maxbytes = 10MB
stderr_logfile = /var/log/newspaper-supervisor-stderr.log
stderr_logfile_maxbytes = 100MB
stderr_logfile_backups = 50
stderr_capture_maxbytes = 10MB
startsecs = 5
stopwaitsecs = 30
stopasgroup = true
